# chef-client.service.conf
[program:chef-client]
command=/usr/bin/chef-client -i 120
